# websockets.ini
[program:websockets]
command=php /var/www/html/erp/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=apache
stdout_logfile=/var/www/html/erp/websockets.log
